# dist/network_exporter.service
[Unit]
Description=Network Exporter
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/network_exporter
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target
